# systemd/minecraft_backup_daily.timer
[Unit]
Description=Save one hourly backup each day

[Timer]
OnCalendar=*-*-* 10:00:00
Persistent=true

[Install]
WantedBy=timers.target
