# systemd/minecraft_backup_hourly.timer
[Unit]
Description=Back up the Minecraft server each hour

[Timer]
OnCalendar=hourly
Persistent=false

[Install]
WantedBy=timers.target
